# rootfs/etc/cron.daily/ntpdate
#!/bin/sh
/usr/sbin/ntpdate -s -pool.ntp.org;(exit 0);
